# config/supervisor.conf
[program:app]
command = python3 /simple-calendar-flask/serice.py 80

[program:nginx]
command = service nginx restart
